# conf/supervisor/cloud_telemetry.conf
[program:telemetry_server]
command=/usr/local/bin/run_cloud_telemetry.sh
#directory=/opt/flask
redirect_stderr=true
startsecs=5
autorestart=true
stdout_logfile=/var/log/telemetry_server.log
