# conf/supervisor/host_agent.conf
[program:telemetry_host_agent]
command=/usr/local/bin/run_host_agent.sh
#directory=/opt/flask
redirect_stderr=true
startsecs=5
autorestart=true
stdout_logfile=/var/log/telemetry_host_agent.log
